# services/cover-watch.timer
[Unit]
Description=astro: run cover-watch every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=false

[Install]
WantedBy=timers.target
