# services/cover-controller.timer
[Unit]
Description=astro: run cover-controller every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=false

[Install]
WantedBy=timers.target
