# services/cover-open-when-dark.service
[Unit]
Description=astro: open lens cover the moment skycam reports dark sky
# Runs on starcam (cover logic lives here). skycam-frames and starcam-frames
# are NFS-mounted from puppy (storage only); the servo is driven by a LOCAL
# cover-move call (bin/cover detects it's on starcam and skips SSH). starcam's
# system python3 has numpy+PIL, so no venv is needed. No act window — fires
# whenever the timer wakes it.
RequiresMountsFor=/home/peter/starcam-frames /home/peter/skycam-frames

[Service]
Type=oneshot
User=peter
Group=peter
Environment=STARCAM_ROOT=/home/peter/starcam-frames
Environment=SKYCAM_MOUNT=/home/peter/skycam-frames
ExecStart=/usr/bin/python3 /home/peter/astro/bin/cover-open-when-dark
StandardOutput=journal
StandardError=journal
