# services/cover-close.timer
[Unit]
Description=astro: close lens cover daily at 07:00 local

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
