# services/plot-brightness.service
[Unit]
Description=astro: render the all-nights brightness plot

[Service]
Type=oneshot
User=peter
ExecStart=/home/peter/astro/services/plot-brightness-run.sh
StandardOutput=journal
StandardError=journal
