# services/publish-astrocam.service
[Unit]
Description=astro: publish yesterday astrocam deliverables to S3

[Service]
Type=oneshot
User=peter
ExecStart=/home/peter/astro-unify/services/publish-astrocam-run.sh
StandardOutput=journal
StandardError=journal
