# services/publish-astrocam.timer
[Unit]
Description=astro: run publish-astrocam daily at 07:00 Europe/London

[Timer]
OnCalendar=*-*-* 07:00:00 Europe/London
Persistent=true

[Install]
WantedBy=timers.target
