# services/storage-report.timer
[Unit]
Description=astro: refresh the storage inventory every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
AccuracySec=15min
Persistent=true

[Install]
WantedBy=timers.target
